# autonomous/detection-requests/t1053_005.yml
technique_id: T1053.005
title: Scheduled Task Creation for Persistence
status: MONITORING
priority: high
mitre_tactic: persistence
requested_by: intel_agent
requested_date: '2026-03-06T00:00:00Z'
intel_report: threat-intel/fawkes/fawkes-ttp-mapping.md
scenario_file: simulator\scenarios\t1053_005.json
sigma_rule: detections\persistence\t1053_005.yml
compiled_lucene: detections\persistence\compiled\t1053_005.lucene
compiled_spl: detections\persistence\compiled\t1053_005.spl
deployed_date: '2026-03-06T23:10:00Z'
elastic_rule_id: cf4b3969-5861-4a0d-bd17-785f829f691b
splunk_saved_search: Scheduled Task Creation for Persistence T1053.005
last_quality_review: '2026-03-18T18:40:41Z'
quality_score: 1.0
fp_rate: 0.0
tp_rate: 1.0
alert_volume_24h: 0
cost_estimate: low
auto_deploy_eligible: true
fawkes_commands:
- schtask
changelog:
- date: '2026-03-06T00:00:00Z'
  agent: intel_agent
  action: created
  details: Fawkes schtask command creates scheduled tasks. Both Fawkes and Scattered Spider use this technique.
- date: '2026-03-06T23:04:09Z'
  agent: blue-team
  action: update
  details: 'Updated fields: sigma_rule, compiled_lucene, compiled_spl'
- date: '2026-03-06T23:04:09Z'
  agent: blue-team
  action: transition:SCENARIO_BUILT->AUTHORED
  details: Sigma rule authored, transpiled to Lucene + SPL
- date: '2026-03-06T23:04:09Z'
  agent: blue-team
  action: update
  details: 'Updated fields: quality_score, fp_rate, tp_rate'
- date: '2026-03-06T23:04:09Z'
  agent: blue-team
  action: transition:AUTHORED->VALIDATED
  details: F1=1.0, FP_rate=0.0, tier=auto_deploy
- date: '2026-03-06T23:04:09Z'
  agent: blue-team
  action: update
  details: 'Updated fields: auto_deploy_eligible'
- date: '2026-03-06T23:10:43Z'
  agent: blue-team
  action: update
  details: 'Updated fields: elastic_rule_id'
- date: '2026-03-06T23:10:43Z'
  agent: blue-team
  action: update
  details: 'Updated fields: splunk_saved_search'
- date: '2026-03-06T23:11:56Z'
  agent: blue-team
  action: transition:VALIDATED->DEPLOYED
  details: Deployed to Elastic + Splunk
- date: '2026-03-06T23:11:56Z'
  agent: blue-team
  action: update
  details: 'Updated fields: deployed_date'
- date: '2026-03-07T23:21:41Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-07T23:21:41Z'
  agent: quality
  action: transition:DEPLOYED->MONITORING
  details: "Health score 0.915 \u2014 healthy"
- date: '2026-03-07T23:34:59Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-07T23:36:41Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-08T14:33:03Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-10T22:34:58Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-11T14:05:00Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-12T18:35:49Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-13T18:25:10Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-14T18:19:57Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-15T18:21:35Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-16T18:43:49Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
- date: '2026-03-18T18:40:41Z'
  agent: quality
  action: update
  details: 'Updated fields: alert_volume_24h, last_quality_review'
