# detections/persistence/t1053_005.yml
title: Scheduled Task Creation for Persistence [T1053.005]
id: f9c8dd06-6eb2-4dec-8c08-f795f9482cc3
status: experimental
description: Detects scheduled task creation for persistence. Detection should flag schtasks /Create with executables in user-writable
  paths (AppData, Temp, ProgramData). Exclude known system/vendor tasks.
author: blue-team-agent
date: '2026-03-06'
modified: '2026-03-06'
references:
- https://attack.mitre.org/techniques/T1053/005/
- https://github.com/galoryber/fawkes
tags:
- attack.persistence
- attack.t1053.005
- detection.fawkes
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    event.code: '1'
    process.name: schtasks.exe
    process.command_line|contains:
      - '/Create'
      - '/create'
  selection_suspicious_target:
    process.command_line|contains:
      - '\AppData\'
      - '\Temp\'
      - '\ProgramData\'
      - '\Users\Public\'
      - '\Downloads\'
  filter_system_tasks:
    process.command_line|contains:
      - '\Windows\System32\defrag'
      - '\Windows\System32\cleanmgr'
      - '\Windows\System32\wsqmcons'
      - 'GoogleUpdate.exe'
      - 'MicrosoftEdgeUpdate'
  condition: selection and selection_suspicious_target and not filter_system_tasks
falsepositives:
  - Legitimate software installers that create scheduled tasks
  - IT management tools deploying via schtasks
  - System administrators creating tasks from non-standard paths
level: high
